# bin/entrypoint.service
[Unit]
Description=Solana Entrypoint
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=solanad
StartLimitIntervalSec=0
LogRateLimitIntervalSec=0
ExecStart=/home/solanad/bin/entrypoint.sh

[Install]
WantedBy=multi-user.target
