# bin/sys-tuner.service
[Unit]
Description=Solana System Tuner
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=sol
StartLimitIntervalSec=0
LogRateLimitIntervalSec=0
ExecStart=/home/sol/bin/sys-tuner.sh

[Install]
WantedBy=multi-user.target
